# roles/borgbackup/files/backup@.service
[Unit]
Description="Backup to repo %I"

[Service]
ExecStart=/usr/local/bin/backup.sh

EnvironmentFile=/usr/local/etc/borg/%i.env
