# roles/nginx-membres/files/systemd/update-ledetour-membres-data.timer
[Unit]
Description=Update members.json file of ledetour-membres every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
