# roles/gdrivesync/files/gdrive-sync.timer
[Unit]
Description=""

[Timer]
OnBootSec=1min
OnUnitActiveSec=24h
Unit=gdrive-sync.service

[Install]
WantedBy=multi-user.target
