# roles/gdrivesync/files/gdrive-sync.service
[Unit]
Description="A script to copy the google drive data to the server"

[Service]
ExecStart=/bin/bash /usr/local/bin/gdrive/gdrive-sync.sh
