# files/borgbackup/backup.service
[Unit]
Description="A script to run borgbackup"

[Service]
ExecStart=/bin/bash /usr/local/bin/backup/backup.sh
